# netscout.service
[Unit]
Description=NetScout Pi API Service
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/NetScout-Pi-V2
# ExecStart will be replaced by setup.sh to use either:
# - The Python virtual environment if created
# - The system Python if not using a virtual environment
ExecStart=/usr/bin/python3 /home/pi/NetScout-Pi-V2/run.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=netscout
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
